# src/systemd/driveraker.timer
[Unit]
Description=Run driveraker every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
