# src/systemd/driveraker.service
[Unit]
Description=Sync driveraker

[Service]
ExecStart=/bin/bash /home/driveraker/.config/driveraker/sync

[Install]
WantedBy=default.target
