# driveraker.service
[Unit]
Description=Sync driveraker

[Service]
Type=oneshot
ExecStart=/bin/bash /home/deleuze/.config/driveraker/sync.sh
